# conf/supervisor.conf
; python后台总开关，可在service/__init__.py 开关具体服务
[program:service]

command     = /startalk/qtalk_search/venv/bin/gunicorn -w6 start_service:app -b 0.0.0.0:8884
environment = PATH="/startalk/qtalk_search/venv/bin",PYTHONPATH="/startalk/qtalk_search:PYTHONPATH"
directory   = service
startsecs   = 3

autostart=true
autorestart=true

redirect_stderr         = true
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups  = 10
stdout_logfile          = ./log/access.log

; 从kafka往redis制作cache的服务
;[program:cache]
;
;command     = /startalk/qtalk_search/venv/bin/python3 /startalk/qtalk_search/service/cache/cache.py
;environment = PATH="/startalk/qtalk_search/venv/bin",PYTHONPATH="/startalk/qtalk_search:PYTHONPATH"
;directory   = service
;startsecs   = 3
;
;autostart=true
;autorestart=true
;
;redirect_stderr         = true
;stdout_logfile_maxbytes = 50MB
;stdout_logfile_backups  = 10
;stdout_logfile          = ./log/cache.log



; 从kafka往es输入数据的服务
;[program:kafka2es]
;
;command     =  /startalk/qtalk_search/venv/bin/python3 /startalk/qtalk_search/service/kafka2es/kafka2es.py
;environment = PATH="/startalk/qtalk_search/venv/bin",PYTHONPATH="/startalk/qtalk_search:PYTHONPATH"
;directory   = service
;startsecs   = 3
;
;autostart=true
;autorestart=true
;
;redirect_stderr         = true
;stdout_logfile_maxbytes = 50MB
;stdout_logfile_backups  = 10
;stdout_logfile          = ./log/2es.log


[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file
;username = jingyu.he
;password = jingyu.he


[supervisord]
logfile=./log/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB         ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10            ; # of main logfile backups; 0 means none, default 10
loglevel=info                 ; log level; default info; others: debug,warn,trace
pidfile=./log/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false                ; start in foreground if true; default false
minfds=1024                   ; min. avail startup file descriptors; default 1024
minprocs=200                  ; min. avail process descriptors;default 200
;umask=022                    ; process file creation umask; default 022
;user=supervisord             ; setuid to this UNIX account at startup; recommended if root
;identifier=supervisor        ; supervisord identifier, default is 'supervisor'
;directory=/tmp               ; default is not to cd during start
;nocleanup=true               ; don't clean up tempfiles at start; default false
;childlogdir=/tmp             ; 'AUTO' child log dir, default $TEMP
;environment=KEY="value"      ; key value pairs to add to environment
;strip_ansi=false             ; strip ansi escape codes in logs; def. false

; The rpcinterface:supervisor section must remain in the config file for
; RPC (supervisorctl/web interface) to work.  Additional interfaces may be
; added by defining them in separate [rpcinterface:x] sections.

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; The supervisorctl section configures how supervisorctl will connect to
; supervisord.  configure it match the settings in either the unix_http_server
; or inet_http_server section.

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[supervisord]
logfile = ./log/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = ./log/supervisord.pid
;nodaemon = true
nodaemon = false
